# xfce4rice5/.config/autostart/picom.desktop
[Desktop Entry]
Name=picom
Comment=An X compositor
Exec=picom --config picom.conf
RunHook=0

